# roles/base/files/etc/systemd/system/wpa_supplicant-nl80211@.service.d/pre.conf
[Service]
ExecStartPre=-/bin/ln -s wpa_supplicant-nl80211.conf /etc/wpa_supplicant/wpa_supplicant-nl80211-%I.conf
