# roles/network/files/etc/systemd/system/wpa_supplicant-nl80211@.service.d/reload.conf
[Unit]
Wants=wpa_supplicant-reload-nl80211@%i.path
