# roles/networking/files/etc/systemd/system/wpa_supplicant-reload-nl80211@.path
[Unit]
Description=WPA supplicant daemon (interface- and nl80211 driver-specific version) configuration file watcher
After=wpa_supplicant-nl80211@%i.service
BindsTo=wpa_supplicant-nl80211@%i.service

[Path]
PathModified=/etc/wpa_supplicant/wpa_supplicant-nl80211-%I.conf
Unit=wpa_supplicant-reload@%i.service
